# deploy/bazi-direction-assistant.service
[Unit]
Description=Bazi Direction Assistant
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/bazi-direction-assistant/current
EnvironmentFile=/opt/bazi-direction-assistant/current/.env.production
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=20
SyslogIdentifier=bazi-direction-assistant

[Install]
WantedBy=multi-user.target
